# systemd/daily-jp.service
[Unit]
Description=Daily Trading Helper JP PLUS (oneshot)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# root 実行（system-level unit）用に絶対パスへ変更
User=root
WorkingDirectory=/home/nishi/InvestmentAlarmForRaspiZero
EnvironmentFile=-/etc/jpplus.env
ExecStart=/home/nishi/InvestmentAlarmForRaspiZero/.venv/bin/python -u -m investment_alarm.daily_system_jp_plus -c config/config.jp.yml
Nice=10
StandardOutput=journal
StandardError=journal
# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
# ProtectHome は /home を遮断するため root 実行では無効化
# （必要なら ReadWritePaths=/home/nishi/InvestmentAlarmForRaspiZero と組み合わせて read-only 運用も可）
# ProtectHome=read-only
# ReadWritePaths=/home/nishi/InvestmentAlarmForRaspiZero
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
UMask=0077
Restart=on-failure

[Install]
WantedBy=default.target
