# systemd/investment-alarms.service
[Unit]
Description=Investment Alarms (US10Y/VIX/USDJPY/SMH-XLK)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/nishi/InvestmentAlarmForRaspiZero
EnvironmentFile=-/etc/jpplus.env
Environment=PYTHONPATH=/home/nishi/InvestmentAlarmForRaspiZero/src
ExecStart=/home/nishi/InvestmentAlarmForRaspiZero/.venv/bin/python -u -m investment_alarm.investment_alarms -c config/config.alarms.yml
Nice=10
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
# ProtectHome は /home を遮断するため root 実行では無効化
# ProtectHome=read-only
# ReadWritePaths=/home/nishi/InvestmentAlarmForRaspiZero
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
UMask=0077
Restart=on-failure

[Install]
WantedBy=default.target
