# systemd/investment-alarms.service
[Unit]
Description=Investment Alarms (US10Y/VIX/USDJPY/SMH-XLK)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/InvestmentAlarmForRaspiZero
EnvironmentFile=-/etc/jpplus.env
ExecStart=%h/InvestmentAlarmForRaspiZero/.venv/bin/python -u investment_alarms.py -c config.alarms.yml
Nice=10
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
UMask=0077
Restart=on-failure

[Install]
WantedBy=default.target
